# scripts/packaging/com.rusty-imap-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.rusty-imap-mcp</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/rusty-imap-mcp</string>
    <string>daemon</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key><false/>
  </dict>
  <key>StandardErrorPath</key><string>/tmp/rusty-imap-mcp.err.log</string>
</dict>
</plist>
